# systemd/coral.service
[Unit]
Description=CORAL Hub - OSINT Monitoring Dashboard
After=network.target

[Service]
Type=simple
User=coral
WorkingDirectory=/opt/coral/coral
ExecStart=/usr/bin/python3 /opt/coral/coral/app.py
Restart=on-failure
RestartSec=10

# Environment
Environment="PYTHONUNBUFFERED=1"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=coral-hub

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
